# gunicorn.service
[Unit]
Description=Gunicorn daemon for Vetgaaf Django application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/vetgaaf
Environment="PATH=/var/www/vetgaaf/venv/bin"
EnvironmentFile=/var/www/vetgaaf/.env
ExecStart=/var/www/vetgaaf/venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/var/www/vetgaaf/gunicorn.sock \
    --timeout 60 \
    --access-logfile /var/www/vetgaaf/logs/gunicorn-access.log \
    --error-logfile /var/www/vetgaaf/logs/gunicorn-error.log \
    rea.wsgi:application

[Install]
WantedBy=multi-user.target
